# scripts/utils/sihsalus-compose.service
[Unit]
Description=SIH Salus Docker Compose stack
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service
ConditionPathExists=/opt/sihsalus/docker-compose.yml

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/sihsalus
EnvironmentFile=-/opt/sihsalus/.env
ExecStartPre=/usr/bin/docker compose config --quiet
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecReload=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose stop
TimeoutStartSec=0
TimeoutStopSec=300
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
